# vorlagen/crontab
* * * * * /home/pi/messungen/skripte/push_all.sh
* * * * * /home/pi/messungen/skripte/check-ip.sh
